# meta-gigcompute/recipes-core/gigcompute-set-hostname/files/gigcompute-set-hostname.service
[Unit]
Description=Set Hostname derived from SOMs serial number
Before=network-pre.target
Wants=network-pre.target
RequiresMountsFor=/etc

# Only execute once on first boot
ConditionPathExists=!/etc/gigcompute_hostname_set

[Service]
Type=oneshot
ExecStart=/usr/sbin/gigcompute-set-hostname.sh

[Install]
WantedBy=sysinit.target
